# holzeinschlag.service
[Unit]
Description=Holzeinschlag Austria Web App
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/holzeinschlag-austria
ExecStart=/home/exedev/holzeinschlag-austria/server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
